# salt/env/base/states/node/config/node-server.conf
{% if 'dev' in grains['host'] %}
{% set environment = 'dev' %}
{% elif 'qa' in grains['host']  %}
{% set environment = 'qa' %}
{% elif 'rc' in grains['host']  %}
{% set environment = 'rc' %}
{% elif 'prod' in grains['host']  %}
{% set environment = 'prod' %}
{% endif %}

{% if 'web' in grains['host'] %}
{% set directory = 'app' %}
{% elif 'api' in grains['host']  %}
{% set directory = 'api' %}
{% endif %}

# node Server Start up script
# this script will start and stop a the node js server

#configuration vagriables
script
cd /var/apps/www/{{ directory }}
NODE_ENV={{ environment }} node server.js
emit node_started
end script

start on startup

